# season-2/platzigram-deploy/roles/platzigram/files/upstart/platzigram.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/platzigram

env NODE_ENV="production"
env PLATZIGRAM_SECRET="d240a6ce-ba49-4ed4-ab35-9c881868d296"
env AWS_ACCESS_KEY="AKIAIFB4IFH7ZFV3OY5A"
env AWS_SECRET_KEY="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"
env FACEBOOK_CLIENT_ID="1216836858350108"
env FACEBOOK_CLIENT_SECRET="ee3b011a73f9a950ab2af671271a3680"

exec npm start
