# season-2/platzigram-deploy/roles/platzigram/files/upstart/platzigram-api-auth.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/platzigram-api

env NODE_ENV="production"
env PLATZIGRAM_SECRET="d240a6ce-ba49-4ed4-ab35-9c881868d296"

exec npm run auth
